# resources/debian/postinst.in
#!/bin/sh

APP=@APP@
PLAY_USER=@USER@
CLUSTER=@CLUSTER@

if test -n "$CLUSTER"
then
	APP_FOLDER=$APP/$CLUSTER
	APP_NAME=$APP-$CLUSTER
else
	APP_FOLDER=$APP
	APP_NAME=$APP
fi

set -e

case "$1" in
configure)
	if grep $PLAY_USER /etc/passwd > /dev/null
	then
		echo User $PLAY_USER already exists
	else
		echo Creating user $PLAY_USER
		adduser --system --home /home/$PLAY_USER --shell /bin/false --no-create-home \
			--group --disabled-password --disabled-login $PLAY_USER
	fi

	chown $PLAY_USER. /var/log/$APP_FOLDER /var/lib/$APP_FOLDER /var/run/$APP_FOLDER

	update-rc.d $APP_NAME defaults 91 09

	echo "Not starting $APP_NAME, configure /etc/$APP_FOLDER/application.conf and then run /etc/init.d/$APP_NAME start"
esac
